# crontabs
43 16 * * * /juejing_actions >> /var/log/cron 2>&1 &